# peter/deploy/systemd/irpete-peter.service
# Template: copy to /etc/systemd/system/irpete-peter.service
# Adjust WorkingDirectory and ExecStart to match your Pi install (venv python).
#
#   sudo cp deploy/systemd/irpete-peter.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now irpete-peter
#
# Secrets and TLS paths live in /etc/irpete/peter.env (chmod 600). See deploy/peter.env.example.

[Unit]
Description=IRPete Peter — HTTPS API (Uvicorn)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# Match WorkingDirectory to the directory that contains `data/` (SQLite default) and your venv.
WorkingDirectory=/opt/irpete/peter
EnvironmentFile=/etc/irpete/peter.env

# Replace with the interpreter from the venv where `pip install -e .` was run.
ExecStart=/opt/irpete/peter/.venv/bin/python -m irpete.main

Restart=always
RestartSec=5

NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
